# scripts/guardian-claude-hw.timer
[Unit]
Description=Guardian Claude daily 7:15am timer for H&W

[Timer]
OnCalendar=*-*-* 07:15:00
Persistent=true

[Install]
WantedBy=timers.target
